# com.analyseai.preload.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.analyseai.preload</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/anayakul/codebase/analyseai/preload-model.sh</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <false/>
    <key>StandardOutPath</key>
    <string>/Users/anayakul/.cache/analyseai/preload.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/anayakul/.cache/analyseai/preload-error.log</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>MLC_DEVICE</key>
        <string>cpu</string>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>
</dict>
</plist> 